# badger/rootfs/etc/s6-overlay/s6-rc.d/badger/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

# Declare variables
declare log_level

## Get the 'message' key from the user config options.
log_level=$(bashio::config 'log_level')
bashio::log.debug "Setting log level to ${log_level:=info}"

# FIXME: Remove dis - very bad for security
cat /data/options.json
ls -alh /data

bashio::log.info "Service is starting..."
set -e
cd /var/run && /usr/bin/badger run --level="${log_level:=info}" --addon-options="/data/options.json"
